# roles/networking/templates/apply-iptables.service.j2
# roles/wifi_access_point/templates/apply-iptables.service.j2
[Unit]
Description=Apply iptables rules for WiFi AP
After=network.target

[Service]
Type=oneshot
ExecStart=/etc/iptables-rules.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target