# roles/power_management/templates/power-monitor.service.j2
[Unit]
Description=Power Monitoring Service for Pi5 Hub
After=network.target docker.service

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /usr/local/bin/power_monitor.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target