# roles/wifi_access_point/templates/frey-wifi-captive-portal-auto.service.j2
[Unit]
Description=Frey Automatic WiFi Connector
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/frey-wifi-captive-portal-auto.sh

[Install]
WantedBy=multi-user.target
