# roles/wifi_access_point/templates/frey-wifi-healthcheck.service.j2
[Unit]
Description=Frey WiFi connectivity watchdog
# Note: NO network-online.target dependency - must work offline
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/frey-wifi-healthcheck
