# roles/wifi_access_point/templates/frey-wifi-captive-portal-auto.timer.j2
[Unit]
Description=Run Frey Automatic WiFi Connector every 2 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
Unit=frey-wifi-captive-portal-auto.service

[Install]
WantedBy=timers.target
